# docker/php/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log

[program:messenger-worker]
command=php /var/www/backend/bin/console messenger:consume async --time-limit=3600
directory=/var/www/backend
autostart=true
autorestart=true
user=www-data
stderr_logfile=/var/log/messenger.err.log
stdout_logfile=/var/log/messenger.out.log